# daemon/systemd/joystick_restart.service
[Unit]
Description=Restarter for joystick virtual keyboard service

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl restart joystick_virt_kbd
